# files/rc.local
#!/bin/bash
/bin/su i2p -c "/home/i2p/i2p/i2prouter start"
exit 0